# launch-agents/com.nocall.voice-input-tool.fixed.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nocall.voiceinputtool.native</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/shogohayashi/Applications/VoiceInputTool.app/Contents/MacOS/VoiceInputTool</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/shogohayashi/voice-input-tool</string>
  <key>StandardOutPath</key>
  <string>/Users/shogohayashi/voice-input-tool/logs/launch-agent.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/shogohayashi/voice-input-tool/logs/launch-agent.err.log</string>
  <key>LimitLoadToSessionType</key>
  <string>Aqua</string>
</dict>
</plist>
